# debian/windmaker-infiniterecorder.cron.d
5 * * * * root /usr/bin/windmaker-infiniterecorder-video-manager
