# debian/windmaker-infiniterecorder.cron.d
*/2 * * * * root flock -x /tmp/windmaker-infiniterecorder-video-manager.lock -c '/usr/bin/windmaker-infiniterecorder-video-manager'
