# debian/windmaker-infiniterecorder.cron.d
0 2 * * * root flock -x -w 30 /tmp/windmaker-infiniterecorder-video-manager.lock -c '/usr/local/bin/windmaker-infiniterecorder-video-manager'
